# quadlet.d/excalidraw.container
[Container]
Image=docker.io/excalidraw/excalidraw:latest
PublishPort=8080:80

[Install]
WantedBy=default.target
