# quadlet.d/hello-world.container
[Container]
Image=docker.io/hello-world:latest

[Install]
WantedBy=default.target