# quadlet.d/hello-world.container
[Container]
Image=docker.io/hello-world:latest

[Service]
StandardError=tty
StandardOutput=tty

[Install]
WantedBy=default.target